# gfs-downloader.service
[Unit]
Description=GFS Weather Data Downloader Daemon
After=network.target mysql.service

[Service]
Type=simple
User=twoj_uzytkownik
WorkingDirectory=/home/twoj_uzytkownik/gfs_downloader
Environment="PATH=/home/twoj_uzytkownik/gfs_downloader/venv/bin"
ExecStart=/home/twoj_uzytkownik/gfs_downloader/venv/bin/python /home/twoj_uzytkownik/gfs_downloader/gfs_downloader_daemon.py
Restart=always
RestartSec=10
StandardOutput=append:/home/twoj_uzytkownik/gfs_downloader/logs/daemon_service.log
StandardError=append:/home/twoj_uzytkownik/gfs_downloader/logs/daemon_service_errors.log

[Install]
WantedBy=multi-user.target

